# config/crontab
1 0 * * * cd <PROJECTPATH>/releases/current/<APP_NAME> && /mnt/ym/<APP_NAME>/bin/python /mnt/ym/<APP_NAME>/releases/current/<APP_NAME>/manage.py run_nightly_tasks >> /var/log/nightly.log 2>&1
